# debian/postrm
#!/bin/sh
# postrm script for igestis
#
# see: dh_installdeb(1)

set -e

if [ -f /usr/share/debconf/confmodule ]; then
    . /usr/share/debconf/confmodule
fi



if [ -f /usr/share/dbconfig-common/dpkg/postrm.mysql ]; then
    . /usr/share/dbconfig-common/dpkg/postrm.mysql
    dbc_go igestis $@
fi





case "$1" in
    purge)
    
	rm -f /usr/share/igestis/config.php
	rm -f /etc/igestis/config.php.dpkg-bak
	rm -f /etc/igestis/config.php
	rm -f /etc/igestis/ConfigIgestisGlobalVars.php
	rm -f /usr/share/igestis/config/igestis/ConfigIgestisGlobalVars.php
	rmdir /usr/share/igestis/ 2>/dev/null || true
	
	rm -f /etc/igestis/debian-db.php
	if which ucf >/dev/null 2>&1; then
		ucf --purge /etc/igestis/debian-db.php
	fi
	
	if [ -e /usr/share/debconf/confmodule ] ; then
		. /usr/share/debconf/confmodule
		db_purge
	fi

	rm -rf /var/cache/igestis/
	rm -rf /var/log/igestis/
	
	ucf --purge /etc/igestis/ConfigIgestisGlobalVars.php
	
    ;;
    
    remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
	db_get igestis/reconfigure-webserver || true
	webservers="$RET"

	for webserver in $webservers; do
		webserver=${webserver%,}
		rm -f /etc/$webserver/conf.d/igestis
		invoke-rc.d $webserver reload 3>/dev/null || true
	done
	
	
	
    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


